# scripts/deploy/systemd/traveler-web.service
[Unit]
Description=Traveler Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/traveler --web --port 8080 --data-dir /home/pi/.traveler
Restart=always
RestartSec=10
Environment=HOME=/home/pi
EnvironmentFile=/home/pi/.traveler/.env

# Logging
StandardOutput=append:/home/pi/.traveler/web.log
StandardError=append:/home/pi/.traveler/web.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/pi/.traveler

[Install]
WantedBy=multi-user.target
